# init/psd-resync.timer
[Unit]
Description=Timer for profile-sync-daemon - 1Hour
PartOf=psd-resync.service psd.service

[Timer]
OnUnitActiveSec=1h
