# psd.cron.hourly
#!/bin/sh
/usr/bin/profile-sync-daemon sync &> /dev/null
