# psd.cron.hourly
#!/bin/sh
/usr/sbin/rc.d sync psd &> /dev/null
