# common/debian/preinst
#!/bin/sh

set -e
set -x
case "$1" in
	install|upgrade)
		# version 6 is a user daemon so clear out old system level files

		[ "$(systemctl is-enabled psd)" = "enabled" ] &&
			systemctl disable psd.service &>/dev/null
		
		[ "$(systemctl is-active psd)" = "active" ] &&
			systemctl stop psd.service &>/dev/null
	
		[ -f /etc/psd.conf ] &&
			mv /etc/psd.conf /etc/psd.conf.old
		;;

	abort-upgrade)
		/bin/true
		;;

	*)
		echo "preinst called with unknown argument \`$1'" >&2
		exit 1
		;;
esac
exit 0
